# src/php-5.3/build/modules/build/files/etc/supervisor/conf.d/php.conf
[program:php]
command = /usr/bin/pidproxy /usr/local/src/phpfarm/inst/php-5.3.29/var/run/php-fpm.pid /usr/local/src/phpfarm/inst/php-5.3.29/sbin/php-fpm -F
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
