# src/build/modules/php/files/etc/supervisor/conf.d/php.conf
[program:php]
command = /phpfarm/inst/php-5.3.29/sbin/php-fpm -F
autostart = true
autorestart = true
