# src/run/modules/php/templates/mysqld.conf.erb
[program:mysqld]
command=socat UNIX-LISTEN:/tmp/mysql.sock,fork,reuseaddr,unlink-early,user=www-data,group=www-data,mode=777 TCP:<%= @mysqld_host %>
autostart=true
autorestart=true
